# crypto_alert_bot.service
[Unit]
Description=Crypto alert bot
StartLimitBurst=5
StartLimitIntervalSec=3600

[Service]
Type=simple
StandardOutput=journal
StandardError=journal
WorkingDirectory=/home/USER/crypto_alert_bot
ExecStart=/home/USER/crypto_alert_bot/run.sh
Restart=always
RestartSec=15
User=USER
Group=USER

[Install]
WantedBy=default.target
